# supervisord.conf
[program:tilestream]
command=/home/dotcloud/node_modules/.bin/tilestream start --uiPort=8080 --tilePort=8080 --host ''

